# rootfs/etc/s6-overlay/s6-rc.d/novnc/run
#!/usr/bin/with-contenv bash

echo "***** starting novnc *****"
/usr/share/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 0.0.0.0:6080
