# rootfs/etc/s6-overlay/s6-rc.d/vnc/run
#!/usr/bin/with-contenv bash

echo "***** starting vnc as $VNC_USERNAME *****"
exec runuser -u $VNC_USERNAME -- vncserver -localhost no -depth 16 -fg :0